# config/systemd/user/dotfiles-backup.service
[Unit]
Description=Dotfiles auto-backup

[Service]
Type=oneshot
ExecStart=%h/dotfiles/backup.sh
